# files/etc/uci-defaults/90-6doll
#!/bin/sh

uci -q batch <<-EOF
	set network.lan.proto='dhcp'
	commit network
EOF
